# ami/packer/ansible/templates/eureka.service.j2
[Unit]
Description=Eureka Service Registry
After=network.target

[Service]
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory=/opt/{{ app_name }}
ExecStart=/usr/bin/java -jar {{ app_name }}.jar --spring.config.location=file:/opt/{{ app_name }}/application.yml
SuccessExitStatus=143
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables (optional - remove if using hardcoded config)
Environment="EUREKA_HOST=%H"
Environment="EUREKA_PORT=8761"
Environment="EUREKA_PREFER_IP=true"
Environment="EUREKA_REGISTER_WITH_EUREKA=false"
Environment="EUREKA_FETCH_REGISTRY=false"

[Install]
WantedBy=multi-user.target
